# systemd/.config/systemd/user/daily-devotional.timer
[Unit]
Description=Run the daily devotional at 5:30 AM

[Timer]
OnCalendar=*-*-* 05:30:00
Persistent=true

[Install]
WantedBy=timers.target
