# systemd/.config/systemd/user/daily-devotional.service
[Unit]
Description=Daily devotional (Claude Code)
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/Scripts/bin/daily-devotional.sh
