# systemd/.config/systemd/user/morning-briefing.service
[Unit]
Description=Daily morning briefing (Claude Code)
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/Scripts/bin/morning-briefing.sh
