# systemd/.config/systemd/user/obsidian-backup.service
[Unit]
Description=Mirror ~/Obsidian (all vaults incl. Becca's) to zeus
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/Scripts/obsidian/obsidian-backup.sh
TimeoutStartSec=900
