# systemd/.config/systemd/user/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway (v2026.3.2)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/node /home/nk/.npm-global/lib/node_modules/openclaw/dist/index.js gateway --port 18789
Restart=always
RestartSec=5
KillMode=process
Environment=HOME=/home/nk
Environment=TMPDIR=/tmp
Environment=PATH=/home/nk/.bun/bin:/home/nk/.local/bin:/home/nk/.npm-global/bin:/home/nk/bin:/home/nk/.volta/bin:/home/nk/.asdf/shims:/home/nk/.nvm/current/bin:/home/nk/.fnm/current/bin:/home/nk/.local/share/pnpm:/usr/local/bin:/usr/bin:/bin
Environment=OPENCLAW_GATEWAY_PORT=18789
Environment=OPENCLAW_GATEWAY_TOKEN=ab6d62e8eebd6b6b1af9734d179617d24a33cfa7bf5da1ca
Environment=OPENCLAW_SYSTEMD_UNIT=openclaw-gateway.service
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment=OPENCLAW_SERVICE_VERSION=2026.3.2

[Install]
WantedBy=default.target
